# src/Soloco.RealTimeWeb/supervisord.conf
[supervisord]
nodaemon=true

[program:kestrel]
directory=/app/MyApp  
command=k kestrel  